# etc/supervisor.d/01_openvpn.ini
[program:openvpn]
command=/usr/sbin/openvpn --cd %(ENV_VPNPATH)s --config %(ENV_VPN)s.ovpn
autorestart=true
priority=10
