# etc/supervisor.d/02_proxy.ini
[program:dante]
command=sockd -N 2 -D
autorestart=true
priority=20
startretries=15
